# alloy/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:alloy]
command=alloy run /etc/alloy/config.alloy --server.http.listen-addr=0.0.0.0:12345
autostart=true
autorestart=true
