# alloy/supervisord.conf
[supervisord]
nodaemon=true

[program:alloy]
command=alloy run /etc/alloy/config.alloy
autostart=true
autorestart=true
